# com.tommy.ripcurl-to-strava.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tommy.ripcurl-to-strava</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/tommy/dev/ripcurl-to-strava/rcs.py</string>
    <string>sync</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tommy/dev/ripcurl-to-strava</string>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/rcs-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/rcs-sync.log</string>
</dict>
</plist>
